# deploy/crypto-telemetry.service
[Unit]
Description=Crypto Collector - VPS Data Collector + Telemetry API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/collector
EnvironmentFile=/opt/collector/.env
Environment=RUN_MODE=collector
ExecStart=/opt/collector/venv/bin/python3.11 -u main.py
Restart=always
RestartSec=5
StandardOutput=append:/opt/collector/logs/collector.log
StandardError=append:/opt/collector/logs/collector.log

[Install]
WantedBy=multi-user.target
