# deploy/crypto-watchdog.service
[Unit]
Description=Crypto Terminal - Telegram Watchdog
After=network.target crypto-collector.service crypto-telemetry.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/collector
EnvironmentFile=/opt/collector/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/collector/venv/bin/python3.11 -u telemetry/watchdog.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/collector/logs/watchdog.log
StandardError=append:/opt/collector/logs/watchdog.log

[Install]
WantedBy=multi-user.target
